# docker/services.d/ragdoll-gateway/run
#!/command/with-contenv bash
exec ragdoll serve
